# os/bolt.service
[Unit]
Description=BOLT 3 Software

[Service]
LimitNICE=39
ExecStart=/home/dash/startup.sh
User=torc

[Install]
WantedBy=multi-user.target
